# new-api.service
[Unit]
Description=New API Service
After=network.target

[Service]
User=gilberthan
WorkingDirectory=/home/gilberthan/disk1/projects/new-api
Environment=GIN_MODE=release
ExecStart=/home/gilberthan/disk1/projects/new-api/new-api --port 2077 --log-dir /home/gilberthan/disk1/projects/new-api/logs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
